# deploy/monitor-bot.service
[Unit]
Description=Monitor Bot — VPS Security Notifications
After=network.target fail2ban.service
Wants=fail2ban.service

[Service]
Type=simple
User=monitor-bot
Group=monitor-bot
WorkingDirectory=/opt/monitor-bot
ExecStart=/opt/monitor-bot/venv/bin/python bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/monitor-bot /opt/monitor-bot

[Install]
WantedBy=multi-user.target
