# deploy/security-bot.service
[Unit]
Description=Security Bot — VPS Security Notifications
After=network.target fail2ban.service
Wants=fail2ban.service

[Service]
Type=simple
User=security-bot
Group=security-bot
WorkingDirectory=/opt/security-bot
ExecStart=/opt/security-bot/venv/bin/python bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/security-bot /opt/security-bot

[Install]
WantedBy=multi-user.target
